# config/roles/gitlab_backup/templates/gitlab-backup.timer.j2
[Unit]
Description=Run GitLab backup to GDrive daily

[Timer]
OnCalendar={{ gitlab_backup_on_calendar }}
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
