# config/roles/vault_server/templates/vault-snapshot.timer.j2
[Unit]
Description=Run Vault Raft snapshot backup every 6 hours

[Timer]
OnCalendar={{ vault_snapshot_on_calendar }}
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
